# root/etc/services.d/tvheadend/run
#!/usr/bin/with-contenv bash

exec \
        s6-setuidgid abc /usr/bin/tvheadend -C -c /config $RUN_OPTS
